# api/laravel-worker.conf
# Laravel Queue Worker Supervisor Configuration
# Production için Supervisor konfigürasyonu
#
# Kurulum:
# 1. Supervisor install edin:
#    sudo apt-get install supervisor
#
# 2. Bu dosyayı /etc/supervisor/conf.d/laravel-worker.conf olarak kopyalayın
#
# 3. Supervisor'u yeniden başlatın:
#    sudo supervisorctl reread
#    sudo supervisorctl update
#    sudo supervisorctl start laravel-worker:*

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/your/api/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/path/to/your/api/storage/logs/worker.log
stopwaitsecs=3600

# Çeşitli queue'lar için farklı worker'lar:
#
# [program:laravel-worker-emails]
# process_name=%(program_name)s_%(process_num)02d
# command=php /path/to/your/api/artisan queue:work --queue=emails --sleep=3 --tries=3
# autostart=true
# autorestart=true
# numprocs=1
#
# [program:laravel-worker-reports]
# process_name=%(program_name)s_%(process_num)02d
# command=php /path/to/your/api/artisan queue:work --queue=reports --sleep=3 --tries=2 --timeout=300
# autostart=true
# autorestart=true
# numprocs=1
#
# [program:laravel-worker-default]
# process_name=%(program_name)s_%(process_num)02d
# command=php /path/to/your/api/artisan queue:work --queue=default --sleep=3 --tries=3
# autostart=true
# autorestart=true
# numprocs=2
